# etc/init.d/discover
#!/bin/sh /etc/rc.common
START=93
EXTRA_COMMANDS=""

start () {
        /usr/sbin/discover&
}

stop () {
        service_stop /usr/sbin/discover
}


